# shipit-api.service
# /etc/systemd/system/shipit-api.service
# Runs the Django API via gunicorn as a systemd service

[Unit]
Description=ShipIt API (Gunicorn)
After=network.target postgresql.service
Wants=postgresql.service

[Service]
User=www-data
Group=www-data
WorkingDirectory=/opt/shipit/backend
EnvironmentFile=/opt/shipit/backend/.env.production
ExecStart=/opt/shipit/backend/.venv/bin/gunicorn shipit_api.wsgi:application -c /opt/shipit/gunicorn.conf.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
